# scripts/oh-my-humanizer-listener.service
[Unit]
Description=oh-my-humanizer Discord Listener
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/devswha/workspace/oh-my-humanizer
# Explicit PATH — lesson from bot-learnings.md 5-day silent failure
Environment=PATH=/home/devswha/.local/bin:/home/devswha/.cargo/bin:/home/devswha/.nvm/versions/node/v22.17.1/bin:/usr/local/bin:/usr/bin:/bin
ExecStartPre=/bin/sh -c 'command -v claude'
ExecStart=/home/devswha/.nvm/versions/node/v22.17.1/bin/node scripts/discord-listener.js
Restart=on-failure
RestartSec=10
TimeoutStopSec=120

[Install]
WantedBy=default.target
