# scripts/patina-component-bridge.service
[Unit]
Description=patina component-only bot bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/devswha/workspace/patina
EnvironmentFile=-/home/devswha/workspace/patina/.env
Environment=PATH=/home/devswha/.local/bin:/home/devswha/.cargo/bin:/home/devswha/.nvm/versions/node/v22.17.1/bin:/usr/local/bin:/usr/bin:/bin
ExecStartPre=/bin/sh -c 'test -n "$PATINA_RUNTIME_CLI" && command -v "$PATINA_RUNTIME_CLI" && command -v node'
ExecStart=/home/devswha/.nvm/versions/node/v22.17.1/bin/node scripts/component-bridge.mjs
Restart=always
RestartSec=5
TimeoutStopSec=30
KillSignal=SIGINT

[Install]
WantedBy=default.target
